# box-scripts/systemd/a11oy-port-guard.service
[Unit]
Description=a-11-oy.com port guard — recover 80/443 if the k3d serverlb grabs them mid-session
After=docker.service network-online.target nginx.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/a11oy-port-guard
