# box-scripts/systemd/box-scripts-drift-check.service
[Unit]
Description=Check the box's host-level helper scripts/units against their committed box-scripts/ copies
Documentation=file:///usr/local/sbin/box-scripts-drift-check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Real push notifications reuse the SAME notifier + credentials file as the
# a11oy-uptime / dns-drift monitors (ntfy / Telegram / Slack-Discord webhook).
# The leading "-" makes the env file optional, so the watcher still runs
# (logging only) before any channel is configured. SMS/email is intentionally
# NOT used here: this box has no Gmail/nodemailer transport and outbound TCP 465
# is firewalled.
EnvironmentFile=-/etc/a11oy-uptime.env
# Distinct ntfy title so box-scripts-drift alerts read differently from the others.
Environment=NOTIFY_TITLE=box-scripts drift
# Always invoke the shared notifier on every drift edge / RECOVERED. It is
# fail-soft: with no channel configured it just exits and the alert is still
# logged, so wiring this is harmless before the env file is filled in.
Environment=NOTIFY_CMD=/usr/local/sbin/a11oy-uptime-notify
ExecStart=/usr/local/sbin/box-scripts-drift-check
