# box-scripts/systemd/istiod-fit-strategy.service
[Unit]
Description=Keep istiod rolling-update strategy single-replica-safe on uds-szl-demo (2-vCPU headroom guard)
Documentation=file:///usr/local/sbin/istiod-fit-strategy
After=docker.service network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/istiod-fit-strategy
