# box-scripts/systemd/istiod-fit-strategy.timer
[Unit]
Description=Periodically ensure istiod's rolling-update strategy + HPA fit the 2-vCPU uds-szl-demo node

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
