# box-scripts/systemd/szl-ns-scratch-stale-watch.service
[Unit]
Description=Alert when a labeled scratch namespace outlives its declared expiry (szl.io/ttl-days) on uds-szl-demo
Documentation=file:///usr/local/sbin/szl-ns-scratch-stale-watch
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/szl-ns-scratch-stale-watch
