# box-scripts/systemd/szl-ns-scratch-watch.timer
[Unit]
Description=Periodically check for untracked scratch namespaces on uds-szl-demo and alert on the edge

[Timer]
OnBootSec=4min
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
