# box-scripts/systemd/szl-box-sync-pull.service
[Unit]
Description=Safely fast-forward the box szl-uds-deployment tree onto origin/main when behind
Documentation=file:///usr/local/sbin/szl-box-sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# `pull` is conservative by construction: it no-ops when EQUAL, fast-forwards
# when BEHIND, resets when DIVERGED-but-no-unique-work, and REFUSES (non-zero,
# logged) when the box is AHEAD or holds unique committed work. It never pushes
# and never discards unique work without --force, so it is safe to run unattended.
ExecStart=/usr/local/sbin/szl-box-sync pull
