# box-scripts/systemd/szl-receipts-retention-watch.service
[Unit]
Description=Runtime watchdog: page if the szl-receipts-retention daily job stops firing/stalls/fails on uds-szl-demo
Documentation=file:///usr/local/sbin/szl-receipts-retention-watch
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Real push notifications use the SAME notifier + credentials file as the
# a11oy-uptime / dns-drift / receipt-chain-watch / szl-receipts-retention
# monitors (ntfy / Telegram / Slack-Discord webhook). The leading "-" makes the
# env file optional, so the watchdog still runs (logging only) before any
# channel is configured.
EnvironmentFile=-/etc/a11oy-uptime.env
# Distinct ntfy title so retention-watchdog alerts read differently from the
# retention job's own alerts and the other monitors.
Environment=NOTIFY_TITLE=szl-receipts retention WATCHDOG
Environment=NOTIFY_CMD=/usr/local/sbin/a11oy-uptime-notify
ExecStart=/usr/local/sbin/szl-receipts-retention-watch
