# box-scripts/systemd/vault-keystore-offbox-backup.timer
[Unit]
Description=Periodically push the latest Vault keystore snapshot off the box to durable encrypted storage (so the receipt signing key survives loss of the box)

[Timer]
# Weekly cadence: the keystore changes only when the Transit signing key is
# re-keyed (rare) or the unseal shares rotate, so a weekly off-box copy is
# plenty while keeping the sensitive ciphertext footprint small. The job
# refreshes the local snapshot first, then ships it; it is a safe no-op when the
# off-box destination is not configured or the cluster is down.
OnBootSec=20min
OnCalendar=weekly
AccuracySec=1h
RandomizedDelaySec=1h
Persistent=true
Unit=vault-keystore-offbox-backup.service

[Install]
WantedBy=timers.target
