# box-scripts/systemd/vault-keystore-offbox-backup.service
[Unit]
Description=Push the latest Vault keystore snapshot (encrypted barrier + unseal shares) OFF the box to durable storage, encrypted at rest, so the receipt signing key survives loss of the box itself
Documentation=file:///usr/local/sbin/vault-keystore-offbox-backup
Documentation=file:///opt/szl/szl-uds-deployment/docs/operations/VAULT_PERSISTENT_RUNBOOK.md
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Shared push notifier + credentials file as the other box monitors
# (a11oy-uptime / dns-drift / receipt-chain-watch). Leading "-" keeps the env
# file optional so the job still runs (log-only) before a channel is configured.
EnvironmentFile=-/etc/a11oy-uptime.env
# The PRIVATE off-box config (encryption recipient + destination). Optional so
# an unprovisioned box runs the job as a safe log-only no-op.
EnvironmentFile=-/etc/vault-keystore-offbox.env
Environment=NOTIFY_TITLE=vault-keystore off-box backup
Environment=NOTIFY_CMD=/usr/local/sbin/a11oy-uptime-notify
ExecStart=/usr/local/sbin/vault-keystore-offbox-backup
